# templates/webhook.service.erb
[Unit]
Description=r10k Webhook Service
After=syslog.target network.target

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/webhook
User=<%= @user %>
ExecStart=/usr/bin/ruby /usr/local/bin/webhook

[Install]
WantedBy=multi-user.target
